# contrib/Docker/conf_files/supervisor/haproxystats-pull.ini
[program:haproxystats-pull]
command = /usr/bin/haproxystats-pull -f /etc/haproxystats/haproxystats.conf
stdout_logfile = /var/log/supervisor/%(program_name)s.log
stderr_logfile = /var/log/supervisor/%(program_name)s.log
autorestart = true